# init/contactsd.service
[Unit]
Description="contacts daemon"

[Service]
ExecStart=/opt/contactsd/contactsd

[Install]
WantedBy=multi-user.target
